# lib/templates/construct.init.d.erb
#!/sbin/runscript

DIR='<%= extra[:construct_dir] %>'

depend() {
  need net
}

start() {
  cd $DIR
  ebegin "Starting to provision using Logical Construct"
  LOGFILE="/root/$(date -u +'%Y-%m-%d;%R').log"
  LOGCON_DEPLOYMENT_PLATFORM=<%= extra[:platform] %> bin/rake -v > $LOGFILE
  eend $?
}

stop() {
  einfo "Stopping Logical Construct [noop]"
}
