# lib/templates/construct.init.d.erb
#!/sbin/runscript

DIR='<%= extra[:construct_dir] %>'

depend() {
  need net
}

start() {
  cd $DIR
  ebegin "Starting to provision using Logical Construct"
  LOGCON_DEPLOYMENT_PLATFORM=<%= extra[:deployment_platform] %> bin/rake
  eend $?
}

stop() {
  einfo "Stopping Logical Construct [noop]"
}
